# infrastructure/quadlet/app.container
# CodeB App Container - Quadlet 설정
# 위치: /etc/containers/systemd/codeb-app.container
# 적용: systemctl --user daemon-reload && systemctl --user enable --now codeb-app
#
# Quadlet은 Podman 컨테이너를 systemd 서비스로 관리합니다.
# - 서버 부팅 시 자동 시작
# - 크래시 시 자동 재시작
# - journalctl로 로그 통합

[Unit]
Description=CodeB Application Container
Documentation=https://github.com/your-org/codeb-server
# PostgreSQL과 Redis가 먼저 시작되어야 함
After=network-online.target codeb-postgres.service codeb-redis.service
Wants=network-online.target
Requires=codeb-postgres.service codeb-redis.service

[Container]
# 이미지 설정
Image=ghcr.io/your-org/codeb-app:latest
ContainerName=codeb-app

# 포트 매핑
PublishPort=3000:3000

# 환경변수
Environment=NODE_ENV=production
Environment=PORT=3000
# 시크릿은 별도 파일로 관리 (환경변수 파일)
EnvironmentFile=/opt/codeb/config/app.env

# 볼륨 마운트
Volume=/opt/codeb/data/app:/app/data:Z
Volume=/opt/codeb/logs/app:/app/logs:Z

# 네트워크
Network=codeb-network

# 헬스체크
HealthCmd=/usr/bin/curl -sf http://localhost:3000/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

# 리소스 제한 (4GB RAM 서버 기준)
# Memory=1G
# CPUQuota=100%

# 보안
ReadOnly=false
NoNewPrivileges=true

[Service]
# 재시작 정책
Restart=always
RestartSec=10s
TimeoutStartSec=300
TimeoutStopSec=30

# 실행 사용자 (rootless)
# User=codeb
# Group=codeb

# 로그
StandardOutput=journal
StandardError=journal

# 환경 변수 (systemd 레벨)
Environment=PODMAN_SYSTEMD_UNIT=%n

[Install]
WantedBy=multi-user.target default.target
